# root/etc/services.d/lsyncd/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc lsyncd -nodaemon -insist /etc/lsyncd.lua